# deploy/rpm/postinstall.sh
systemctl daemon-reload
systemctl enable --now crispy-dummy.service
